# roles/consul/templates/consul.service
[Unit]
Description=Consul Service

[Service]
ExecStart=consul agent -dev

[Install]
WantedBy=multi-user.target